# templates/consul-online.target.j2
[Unit]
Description=Consul Online
RefuseManualStart=true
Requires=consul-online.service consul.service
After=consul-online.service consul.service
